# launchd/reminderd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>org.reminderd</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>python3</string>
  <string>__HOME__/Reminderd/src/reminderd</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>SessionCreate</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
      <key>DISPLAY</key>
      <string>:0</string>
      <key>USER</key>
      <string>__USER__</string>
      <key>HOME</key>
      <string>__HOME__</string>
    </dict>
    <key>WorkingDirectory</key>
  <string>__HOME__/Reminderd</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/reminderd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/reminderd.err.log</string>
  </dict>
</plist>
